# extra/cups/init.d-cupsd
#!/sbin/runscript

name="CUPS Printing Service"
command="/usr/sbin/cupsd"
command_args=""

depend() {
	use net
	need dbus
	before nfs
	after logger
}
