# extra/openldap/slapd
#!/sbin/runscript

# FIXME: THIS FILE IS INCOMPLETE
# PLEASE CHECK slapd.in AND CREATE A PROPER OpenRC INIT SCRIPT

depend() {
    # FIXME
}

start() {
    ebegin "Starting OpenLDAP"
    if [ -z "$SLAPD_SERVICES" ]; then
      start-stop-daemon --start --pidfile /var/run/openldap/slapd.pid --exec \
        /usr/sbin/slapd -u ldap -g ldap $SLAPD_OPTIONS
    else
      start-stop-daemon --start --pidfile /var/run/openldap/slapd.pid --exec \
        /usr/sbin/slapd -u ldap -g ldap -h "$SLAPD_SERVICES" $SLAPD_OPTIONS
    fi
    eend $?
}

stop() {
    ebegin "Stopping OpenLDAP"
    start-stop-daemon --stop --pidfile /var/run/openldap/slapd.pid
    retval=$?

    eend ${retval}

    rm -rf /var/db/slapd.pid
    rm -rf /var/db/slapd.args

    return ${retval}
}
