# community/libvirt/libvirtd.in
#
# PROVIDE: libvirtd
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# Add the following line to /etc/rc.conf[.local] to enable libvirtd(8)
#
# libvirtd_enable (bool):       Set to "NO" by default.
#                               Set to "YES" to enable libvirtd.
#

. /etc/rc.subr

name="libvirtd"
rcvar=libvirtd_enable

command="/usr/sbin/libvirtd"
command_args="--daemon"

load_rc_config $name

: ${libvirtd_enable="NO"}

run_rc_command "$1"
