# community/synergy/init.d-synergy
#!/sbin/runscript

depend() {
	use logger dns
	need net
}

start() {
	ebegin "Starting ${SVCNAME}"
	start-stop-daemon --start --exec /usr/local/bin/synergys \
	    --pidfile /var/run/synergys.pid \
	    -- --config /usr/local/etc/synergy.conf
	eend $?
}

stop() {
	ebegin "Stopping ${SVCNAME}"
	start-stop-daemon --stop --exec /usr/local/bin/synergys \
	    --pidfile /var/run/synergys.pid
	eend $?
}
