# extra/slim/init.d-slim
#!/sbin/runscript
name="slim"
command=/usr/bin/slim

depend() {
  need localmount dbus
  after bootmisc
  keyword -jail -prefix
}

start()
{
	ebegin "Starting Simple Login Manager"
	start-stop-daemon --background --start --exec ${command}
	eend $?
}

stop() {
	ebegin "Stopping Simple Login Manager"
	start-stop-daemon --stop --exec ${command}
	eend $?
}
