# extra/slim/slim
#!/sbin/runscript
# Copyright 1999-2011 Gentoo Foundation
# Distributed under the terms of the GNU General Public License, v2 or later

extra_started_commands="reload"

depend() {
    need localmount
    after dbus
}

start() {
    ebegin "Starting Slim Login Manager"

    # We need to test if /var/run/slim exists, since script will fail if it does not
    [ ! -e /var/run/slim ] && mkdir /var/run/slim

    start-stop-daemon --start --pidfile /var/run/slim/pid --exec /usr/bin/slim -- -d

    eend $?
}

stop() {
    local retval

    ebegin "Stopping Slim Login Manager"

    start-stop-daemon --stop --pidfile /run/slim/pid
    retval=$?

    eend ${retval}

}

