# extra/hal/hald
#!/sbin/runscript
# Copyright 1999-2011 Gentoo Foundation
# Distributed under the terms of the GNU General Public License, v2 or later


depend() {
	need dbus
	after bootmisc
}

start() {
	ebegin "Starting Hardware Abstraction Layer daemon"
	start-stop-daemon --start -q --exec /usr/sbin/hald --
	eend $?
}

stop() {
	ebegin "Stopping Hardware Abstraction Layer daemon"
	start-stop-daemon --stop -q --pidfile /var/run/hald/hald.pid 
	eend $?
}
