# extra/mpd/musicpd.in
#!/bin/sh

# $FreeBSD$
#
# PROVIDE: mpd
# REQUIRE: mixer LOGIN avahi_daemon oss
# KEYWORD: shutdown

# Add the following line to /etc/rc.conf to enable mpd:
#
# mpd_enable=YES

. /etc/rc.subr

name=mpd
rcvar=mpd_enable

load_rc_config $name

: ${mpd_enable:=NO}

command=/usr/bin/mpd
command_args=/etc/mpd.conf
required_files=/etc/mpd.conf
start_precmd=${name}_getpidfile
stop_precmd=${name}_getpidfile

mpd_getpidfile()
{
	if get_pidfile_from_conf pid_file /etc/mpd.conf ; then
		pidfile="$_pidfile_from_conf"
	else
		pidfile="/var/run/mpd.pid"
	fi
}

run_rc_command $1
