# community/noip/noip
#!/sbin/runscript
# Copyright 1999-2010 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/app-admin/supervisor/files/init.d,v 1.1 2010/03/09 13:45:32 djc Exp $

DAEMON="/usr/bin/noip2"
DAEMON_OPTS="-c /etc/no-ip2.conf"
checkconfig() {
	if [ ! -f /etc/no-ip2.conf ] ; then
		eerror "Please create /etc/no-ip2.conf:"
		return 1
	fi
	return 0
}

start() {
	checkconfig || return $?
	ebegin "Starting supervisord"
	start-stop-daemon --start \
		--exec ${DAEMON} -b --pidfile /var/run/noip2.pid -- ${DAEMON_OPTS}
	eend $?
}

stop() {
	ebegin "Stopping supervisord"
	start-stop-daemon --stop --pidfile /var/run/noip2.pid
	eend $?
}
